# systemd-user-unit-files/.config/systemd/user/caffeine.service
[Unit]
Description=Prevent screensaver and display sleep

[Service]
ExecStart=/usr/bin/caffeine
Restart=always

[Install]
WantedBy=default.target
