# systemd-user-unit-files/.config/systemd/user/keepassxc.service
[Unit]
Description=Password manager

[Service]
ExecStart=/usr/bin/keepassxc
Restart=always

[Install]
WantedBy=default.target
