# systemd-user-unit-files/.config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/nm-applet
Restart=always

[Install]
WantedBy=default.target
