# systemd-user-unit-files/.config/systemd/user/pasystray.service
[Unit]
Description=System tray applet for PulseAudio

[Service]
ExecStart=/usr/bin/pasystray
Restart=always

[Install]
WantedBy=default.target
