# systemd-user-unit-files/.config/systemd/user/picom.service
[Unit]
Description=Standalone compositor for Xorg

[Service]
ExecStart=/usr/bin/picom
Restart=always

[Install]
WantedBy=default.target
