# systemd-user-unit-files/.config/systemd/user/udiskie.service
[Unit]
Description=Automounter for removable media

[Service]
ExecStart=/usr/bin/udiskie --tray
Restart=always

[Install]
WantedBy=default.target
